# home/dotfiles/autostart/plasma-tahoe-layout.desktop
[Desktop Entry]
Type=Application
Name=Plasma Tahoe Layout
Exec=/home/weshy/.local/bin/plasma-tahoe-layout --once
X-KDE-autostart-after=panel
X-KDE-autostart-phase=2
OnlyShowIn=KDE;
